# etc/cron.d/pihole_chmod_chown
* * * * * root /bin/sh /var/www/html/pihole/reset_database.sh > /var/log/chmod_chown.log